# udev/70-venator.rules
# SPDX-License-Identifier: GPL-2.0-only
#
# Grant the 'venator' group write access to the LED control attributes
# created by the venator kernel module, so userspace tooling
# (cli/venator, future GUI) doesn't need sudo.
#
# Install (handled by `make install`):
#   - copy this file to /etc/udev/rules.d/70-venator.rules
#   - sudo groupadd -r venator
#   - sudo usermod -aG venator $USER
#   - sudo udevadm control --reload
#   - sudo udevadm trigger -s venator
#   - newgrp venator    (or log out / back in for the new group to take)

ACTION=="add", SUBSYSTEM=="venator", \
    RUN+="/bin/sh -c '\
        cd /sys$env{DEVPATH} 2>/dev/null && \
        chgrp venator mode color brightness effect_id apply frame \
            speed direction set modes \
            health_mode calibration_mode charge_control_end_threshold \
            lcd_overdrive boot_sound \
            2>/dev/null; \
        chmod g+w mode color brightness effect_id apply frame \
            speed direction set modes \
            health_mode calibration_mode charge_control_end_threshold \
            lcd_overdrive boot_sound \
            2>/dev/null'"
